# BookClub-final_branch/bookclub/gunicorn.service
[Unit]
Description=Gunicorn instance for a simple hello world app
After=network.target
[Service]
User=ubuntu
Group=www-data
WorkingDirectory=/home/ubuntu/bookclub-docker
ExecStart=/home/ubuntu/bookclub-docker/venv/bin/gunicorn -b localhost:8000 app:app
Restart=always
[Install]
WantedBy=multi-user.target
